# vm-setup/vm-autostart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ansible-nas.vm.autostart</string>
    
    <key>ProgramArguments</key>
    <array>
        <!-- For UTM -->
        <!--
        <string>/Applications/UTM.app/Contents/MacOS/utmctl</string>
        <string>start</string>
        <string>Ubuntu Server</string>
        <string>- -headless</string>
        -->
        
        <!-- For VMware Fusion -->
        <!--
        <string>/Applications/VMware Fusion.app/Contents/Library/vmrun</string>
        <string>start</string>
        <string>/Users/mk/Virtual Machines/Ubuntu-NAS.vmwarevm/Ubuntu-NAS.vmx</string>
        <string>nogui</string>
        -->
        
        <!-- For VirtualBox -->
        <!--
        <string>/usr/local/bin/VBoxManage</string>
        <string>startvm</string>
        <string>Ubuntu-NAS</string>
        <string>- -type</string>
        <string>headless</string>
        -->
        
        <!-- Replace with your actual VM software and paths -->
        <string>/bin/sh</string>
        <string>-c</string>
        <string>echo "Configure this file with your VM software" | logger -t ansible-nas-vm</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>StandardOutPath</key>
    <string>/tmp/ansible-nas-vm.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/ansible-nas-vm.error.log</string>
    
    <key>StartInterval</key>
    <integer>300</integer>
</dict>
</plist>